# DEBIAN/postrm
#!/bin/bash

set -e

echo "Disabling background music daemon..."
systemctl disable bgm
systemctl stop bgm
systemctl daemon-reload

#DEBHELPER#